# loops/systemd/hermes-loop@.service
# hermes-loop@<pack>.service — ein Nacht-Lauf des Loop-Runner-Packs <pack>.
# Installation: loops/systemd/install.sh · Start: systemctl --user start hermes-loop@<pack>
# Der Runner selbst hält Pack- und Repo-Locks; Doppel-Starts enden sauber (rc 3).
[Unit]
Description=Hermes Loop-Runner: Pack %i (night)

[Service]
Type=oneshot
# Läuft im Live-Checkout; der Runner arbeitet in einem isolierten Worktree
# unter ~/.hermes/loops/%i/wt. Standard-Packs pushen/deployen nie; die einzige
# code-/hash-gebundene Ausnahme `dashboard-experience` darf nach Opus-Verifier-PASS und
# Runner-Gates ff-only nach piet-fork landen. Deploy bleibt immer ausgeschlossen.
# Über `bash` gestartet (nicht direkt), damit ein fehlendes x-Bit am shim NIE
# wieder ein 203/EXEC-Silent-Fail gibt (2026-07-03 UI-Start-Bug).
ExecStart=/bin/bash %h/.hermes/hermes-agent/loops/shim.sh %i
# Nacht-Budget: plan + run (max_hours default 7h) + Puffer.
TimeoutStartSec=9h
Nice=10
